# com.xauusd.signals.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.xauusd.signals</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/omar/Desktop/ML/xauusd_signals/start_signal_engine.sh</string>
        <string>start</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/omar/Desktop/ML/xauusd_signals/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/omar/Desktop/ML/xauusd_signals/launchd.error.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/omar/Desktop/ML/xauusd_signals</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>

